# systemd/cookie-guard-spoa.service
[Unit]
Description=Cookie Guard SPOA (haproxy service)
After=network-online.target
Wants=network-online.target

[Service]
# Load CLI options from /etc/default/cookie-guard-spoa if present.
EnvironmentFile=-/etc/default/cookie-guard-spoa
# Listen ONLY on localhost; no TLS, plain TCP. Optional flags can be toggled
# in the default file via COOKIE_GUARD_FLAG_* variables.
ExecStart=/usr/local/bin/cookie-guard-spoa \
  $COOKIE_GUARD_SPOA_OPTS \
  $COOKIE_GUARD_FLAG_DEBUG \
  $COOKIE_GUARD_FLAG_COOKIE_SECURE \
  $COOKIE_GUARD_FLAG_NO_UA_BIND \
  $COOKIE_GUARD_FLAG_ALTCHA_DISABLE \
  $COOKIE_GUARD_FLAG_EXTRA

# Run unprivileged
User=haproxy
Group=haproxy

Restart=on-failure
RestartSec=2
WorkingDirectory=/etc/cookie-guard-spoa
Environment=GOTRACEBACK=crash

# Security hardening (safe with localhost-only)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectClock=yes
ProtectHostname=yes

[Install]
WantedBy=multi-user.target
